# system-boot/rpi-ansible/rpi-setup.service
[Unit]
After=ssh.service

[Service]
ExecStart=/home/ubuntu/ansible/rpi-setup.sh

[Install]
WantedBy=default.target
